# systemd/claude-work.service
[Unit]
Description=ClaudeWork - Claude Code Session Manager
Documentation=https://github.com/windschord/claude-work
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude-work
Group=claude-work
WorkingDirectory=/opt/claude-work
EnvironmentFile=/etc/claude-work/env
# npx github:windschord/claude-work でフォアグラウンド起動
# ネットワーク要件:
#   - 初回実行時は必ず GitHub への接続が必要（パッケージ取得）
#   - 2回目以降は npm キャッシュを利用するが、キャッシュ有効期限などの設定により
#     npx/npm が定期的に更新確認やダウンロードのためネットワークアクセスを行う可能性あり
# 初回実行時に自動セットアップ（Prisma、DB、ビルド）を実行
# HOME を設定する理由:
#   1. npm キャッシュを /opt/claude-work/.npm に配置（HOME=/opt/claude-work 配下に保持）
#   2. ProtectHome=read-only との競合を回避
#   3. 自己完結型インストールを維持
Environment=HOME=/opt/claude-work
ExecStart=/usr/bin/npx github:windschord/claude-work
# 初回起動時のビルドに時間がかかるためタイムアウトを延長
TimeoutStartSec=300
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-work

# セキュリティ設定
NoNewPrivileges=true
# strict: ファイルシステム全体を読み取り専用に（ReadWritePaths で許可したパス以外）
ProtectSystem=strict
# ProtectHome=read-only: ホームディレクトリへの書き込みを制限
# /home 配下のプロジェクトを使用する場合は docs/SYSTEMD_SETUP.md の
# 「セキュリティ設定による制限事項」を参照してください
ProtectHome=read-only
PrivateTmp=true
# アプリケーションディレクトリへの書き込みを許可（データ、キャッシュ、ビルド出力）
ReadWritePaths=/opt/claude-work

[Install]
WantedBy=multi-user.target
