# systemd/claude-work.service
[Unit]
Description=ClaudeWork - Claude Code Session Manager
Documentation=https://github.com/windschord/claude-work
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude-work
Group=claude-work
WorkingDirectory=/opt/claude-work
EnvironmentFile=/etc/claude-work/env
# npx ではなくビルド済みの CLI を直接実行（セキュリティ向上）
ExecStart=/usr/bin/node /opt/claude-work/dist/src/bin/cli.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-work

# セキュリティ設定
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/opt/claude-work/data

[Install]
WantedBy=multi-user.target
